# packaging/systemd/irlume-reconcile.service
[Unit]
Description=Re-apply irlume face-auth PAM wiring after a distro update
Documentation=https://github.com/archledger/irlume
# Runs from TWO triggers: the irlume-reconcile.path unit (a greeter PAM file
# changed) AND at boot (WantedBy below), so a strip that landed while the path
# watcher was inactive (an offline dnf/apt upgrade, or before paths.target on a
# first install) is still caught on the next boot. No ConditionPathExists: the
# command self-gates. With no marker it either ADOPTS an existing wiring into a
# marker (an upgrade from a pre-marker version) or exits; with a marker it
# re-adds only the greeter lines that went missing. Idempotent; the password
# stays the fallback throughout.

[Service]
Type=oneshot
ExecStart=/usr/bin/irlume login reconcile

[Install]
WantedBy=multi-user.target
